# Курсовая работа на профессии "DevOps-инженер с нуля" /ansible/nginx_exporter/templates/prometheus_nginxlog_exporter.service.j2
[Unit]
Description=Nginx Log Exporter

[Service]
User=admin
ExecStart=/opt/prometheus_nginxlog_exporter/prometheus-nginxlog-exporter --log-file=/var/log/nginx/access.log 

[Install]
WantedBy=multi-user.target